# deploy/cocoro-studio.service
[Unit]
Description=COCORO Studio (Next.js 16 フロントエンド)
After=network.target cocoro-influencer.service
Wants=cocoro-influencer.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/cocoro-influencer/frontend

# API バックエンドのURL (FastAPI)
Environment=NODE_ENV=production
Environment=API_ORIGIN=http://localhost:8082
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0

ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
